# roles/dual-screen/templates/90-keyboard-hotplug.rules
ACTION=="add", SUBSYSTEM=="usb", ENV{ID_MODEL_ID}="0141", RUN+="/usr/local/bin/keyboard-hotplug"
